# rc.d/btcd
#!/bin/sh
#
# $OpenBSD$

daemon_user="_btc"
daemon="/home/todd/go/bin/btcd"
daemon_flags="-C ~${daemon_user}/.btcd/btcd.conf"

. /etc/rc.d/rc.subr

pexp="btcd.*btcd.conf"

rc_reload=NO

rc_start() {
	${rcexec} "/usr/bin/time -l ${daemon} ${daemon_flags}" 2>&1 | \
		btclog main $$ &
	renice -n -15 $(pgrep btcd)
}

rc_cmd $1
